# systemd/user-midnighttpd.service.in
[Unit]
Description=The midnighttpd web server (user mode)
Documentation=man:midnighttpd(8)

[Service]
Type=simple
ExecStart=%bindir%/midnighttpd ${HOME}/http/ -u %rundir%/http/${USER}.socket -c ${HOME}/.config/midnighttpd.conf
Restart=on-failure

[Install]
WantedBy=default.target
